# supervisord.conf
[supervisord]
nodaemon=true

[program:snapshotter]
command=/opt/cassandra/bin/snapshotScheduler.sh


[program:cassandra]
command=/opt/cassandra/bin/stups-cassandra.sh