# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:cassandra]
command=/opt/cassandra/bin/stups-cassandra.sh